# net-mgmt/p5-FusionInventory-Agent/files/fusioninventory.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: fusioninventory
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# fusioninventory_enable (bool):	Set to NO by default.
#					Set to YES to enable lldpd.

. /etc/rc.subr

name=fusioninventory
rcvar=fusioninventory_enable

load_rc_config $name

: ${fusioninventory_enable:="NO"}

command=%%PREFIX%%/bin/${name}-agent
pidfile=/var/run/${name}/${name}.pid

command_args="-d"

run_rc_command "$1"
